# run/supervisor-mall.conf
[program:gunicorn]
command= bash /home/hike/mall/gunicorn.sh
directory= /code/mall
user=hike
autostart=true
autorestart=true
redirect_stderr=true
